# apps/ios/scripts/com.simonri.bessel.refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Weekly free-provisioning refresh: rebuilds and reinstalls both the dev and
  prod Bessel apps on the iPhone Mondays and Fridays at 09:00, so neither's
  7-day signature ever expires (max gap between runs is ~4 days). launchd
  runs missed jobs at next wake.

  Install:
    cp apps/ios/scripts/com.simonri.bessel.refresh.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.simonri.bessel.refresh.plist
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.simonri.bessel.refresh</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/simon/dev/bessel/apps/ios/scripts/deploy.sh</string>
  </array>
  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Weekday</key><integer>1</integer>
      <key>Hour</key><integer>9</integer>
      <key>Minute</key><integer>0</integer>
    </dict>
    <dict>
      <key>Weekday</key><integer>5</integer>
      <key>Hour</key><integer>9</integer>
      <key>Minute</key><integer>0</integer>
    </dict>
  </array>
  <key>StandardOutPath</key>
  <string>/tmp/bessel-ios-refresh.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/bessel-ios-refresh.log</string>
</dict>
</plist>
